# service/second/run
#!/bin/sh

# Redirect stderr to stdout
exec 2>&1

echo "Starting the second service"

# An endless loop to simulate a running server
exec chpst -u second sh -c "id; while true; do sleep 100; done"
